# pi/systemd/soltaro-qendercore.service
[Unit]
Description=Soltaro Qendercore to Hubitat bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/soltaro-qendercore
ExecStart=/usr/bin/python3 /home/pi/soltaro-qendercore/soltaro_sync.py daemon --env /home/pi/soltaro-qendercore/soltaro.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
