# pi/systemd/soltaro-qendercore.service
[Unit]
Description=Soltaro Qendercore to Hubitat bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jarvis
WorkingDirectory=/home/jarvis/soltaro-qendercore
ExecStart=/usr/bin/python3 /home/jarvis/soltaro-qendercore/soltaro_sync.py daemon --env /home/jarvis/soltaro-qendercore/soltaro.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
